# etc/profile.d/maven.sh
export MAVEN_HOME=/opt/maven/latest
export PATH=$MAVEN_HOME/bin:$PATH
