# etc/profile.d/sbt.sh
export SBT_HOME=/opt/sbt
export PATH=$SBT_HOME:$PATH
