# etc/profile.d/sencha.sh
export SENCHA_SDK_HOME=/opt/sencha/sdk/latest
export PATH=$SENCHA_SDK_HOME:$SENCHA_SDK_HOME/appbuilder:$SENCHA_SDK_HOME/command:$SENCHA_SDK_HOME/jsbuilder:$PATH
